# ansible/files/btrfs-maint/btrfs-maint.service
[Unit]
Description=Runs btrfs volumes maintenance

[Service]
User=root
Group=root
Type=simple
ExecStart=/usr/local/bin/btrfs-maint
